# igor47.dc.service
[Unit]
Description=igor's docker services
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/home/igor47/repos/igor.services
Type=oneshot
RemainAfterExit=true
ExecStart=docker compose up --wait --detach --remove-orphans
ExecReload=docker compose up --wait --detach --remove-orphans
ExecStop=docker compose down

[Install]
WantedBy=multi-user.target
